# settings/wallpaper/.config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper System

[Service]
Type=oneshot
ExecStart=/home/kkokou/.config/i3/scripts/wallpapper_service.sh

[Install]
WantedBy=graphical.target
