# settings/wallpaper/.config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper System

[Service]
Type=oneshot
ExecStart=/usr/bin/wallpapper.sh service

[Install]
WantedBy=graphical.target
